# sv/xvfb/run
#!/bin/bash

#sv start mysql || exit 1

source /etc/envvars
env

exec 2>&1
exec Xvfb :99 -screen 1 1024x768x24 -ac +extension RANDR
